# linchpin/roles/schedule/templates/kstests-cloud.service.j2
[Unit]
Description=Run kickstart tests in cloud

[Service]
Type=simple
ExecStart={{ run_script_dir_path }}/{{ run_script_name }}
